# automations/vesra-outbound/deploy/systemd/vesra-daily-orchestrator.service
[Unit]
Description=Vesra daily contact lifecycle orchestrator
Wants=network-online.target
After=network-online.target vesra-daily-enrichment.service

[Service]
Type=oneshot
WorkingDirectory=/opt/nightfall/automations/vesra-outbound
EnvironmentFile=/etc/vesra/outbound.env
ExecStartPre=/usr/bin/mkdir -p /var/log/vesra /var/lib/vesra/orchestration-runs
ExecStart=/bin/bash -lc 'set -euo pipefail; python3 lead-gen/scripts/orchestration/run_daily_orchestrator.py --rebuild-queue --summary'
StandardOutput=append:/var/log/vesra/daily-orchestrator.log
StandardError=append:/var/log/vesra/daily-orchestrator.log
