# automations/outbound/systemd/outbound-agent-send-approved@.service
[Unit]
Description=Nightfall outbound approved sender for tenant %i
Wants=network-online.target
After=network-online.target outbound-daily-orchestrator@%i.service

[Service]
Type=oneshot
WorkingDirectory=/opt/nightfall/automations/outbound
Environment=NIGHTFALL_OUTBOUND_TENANT=%i
EnvironmentFile=-/etc/nightfall/outbound/%i.env
ExecStartPre=/usr/bin/mkdir -p /var/log/nightfall/outbound
ExecStart=/bin/bash -lc 'set -euo pipefail; .venv/bin/python tools/run_tenant_tool.py --tenant "%i" send-approved --limit "${OUTBOUND_AGENT_SEND_LIMIT:-5}" --send'
StandardOutput=append:/var/log/nightfall/outbound/%i-send-approved.log
StandardError=append:/var/log/nightfall/outbound/%i-send-approved.log
