# automations/outbound/deploy/systemd/vesra-daily-enrichment.service
[Unit]
Description=Vesra daily enrichment
Wants=network-online.target
After=network-online.target vesra-daily-discovery.service

[Service]
Type=oneshot
WorkingDirectory=/opt/nightfall/automations/outbound
EnvironmentFile=/etc/vesra/outbound.env
ExecStartPre=/usr/bin/mkdir -p /var/log/vesra
ExecStart=/bin/bash -lc 'set -euo pipefail; .venv/bin/python lead-gen/scripts/enrichment/apply_known_enrichment.py; .venv/bin/python lead-gen/scripts/enrichment/enrich_decision_makers.py; .venv/bin/python lead-gen/scripts/enrichment/enrich_public_web.py --limit "${VESRA_DAILY_ENRICH_LIMIT:-50}" --max-pages "${VESRA_DAILY_ENRICH_MAX_PAGES:-2}" --checkpoint'
StandardOutput=append:/var/log/vesra/daily-enrichment.log
StandardError=append:/var/log/vesra/daily-enrichment.log
