# automations/outbound/deploy/systemd/vesra-agent-send-approved.timer
[Unit]
Description=Run Vesra approved email sender every UK weekday

[Timer]
OnCalendar=Mon..Fri *-*-* 10:15:00 Europe/London
Persistent=true
RandomizedDelaySec=10m
Unit=vesra-agent-send-approved.service

[Install]
WantedBy=timers.target
