# automations/outbound/deploy/systemd/vesra-outbound-webhook.service
[Unit]
Description=Vesra outbound unsubscribe and reply webhook
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/nightfall/automations/outbound
EnvironmentFile=/etc/vesra/outbound.env
ExecStart=/opt/nightfall/automations/outbound/.venv/bin/python lead-gen/scripts/outreach/unsubscribe_server.py --host 127.0.0.1 --port 8088
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
