# automations/outbound/deploy/systemd/vesra-daily-orchestrator.service
[Unit]
Description=Vesra daily agentic contact orchestrator
Wants=network-online.target
After=network-online.target vesra-daily-enrichment.service

[Service]
Type=oneshot
WorkingDirectory=/opt/nightfall/automations/outbound
EnvironmentFile=/etc/vesra/outbound.env
ExecStartPre=/usr/bin/mkdir -p /var/log/vesra /var/lib/vesra/orchestration-runs
ExecStart=/bin/bash -lc 'set -euo pipefail; .venv/bin/python lead-gen/scripts/orchestration/run_agentic_orchestrator.py --rebuild-queue --execute-safe-tools --max-tool-runs "${VESRA_AGENTIC_MAX_TOOL_RUNS:-3}" --summary'
StandardOutput=append:/var/log/vesra/daily-orchestrator.log
StandardError=append:/var/log/vesra/daily-orchestrator.log
